# deploy/production/supervisor-queue.conf
[program:panificadora-queue]
command=php /var/www/panificadora/artisan queue:work redis --sleep=3 --tries=3 --timeout=60 --memory=128
user=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/panificadora-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
